# config-linux/.config/systemd/user/hyprmoncfgd.service.d/override.conf
[Service]
ExecStart=
ExecStartPre=/usr/bin/sleep 12
ExecStart=/usr/bin/hyprmoncfgd --debounce 2s --poll-interval 8s
